# misc/dspam-milter.conf.upstart
# Upstart init script for dspam-milter

author "Tom Hendrikx <pydspam@whyscream.net>"
description "dspam milter interface"

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn

setuid dspam
setgid dspam
exec /usr/bin/dspam-milter
